# self-check.timer
[Unit]
Description=Run system self-check every 15 minutes
Requires=self-check.service

[Timer]
OnCalendar=*:0/5
RandomizedDelaySec=2m
Persistent=true

[Install]
WantedBy=timers.target